# roles/userspace/templates/systemd.service.j2
[Unit]
Description={{ task.description }}

[Service]
Type=simple
ExecStart={{ task.execute }}

[Install]
WantedBy=default.target
